# thinkmachine/electron/.config/profile.d/electron.sh
# XDG dir spec compliance
export ELECTRON_CACHE="${XDG_CACHE_HOME:-$HOME/.cache}"/electron
