# thinkmachine/X11/.config/profile.d/xorg-auth.sh
# XDG dir spec compliance
export XAUTHORITY="${XDG_DATA_HOME:-$HOME/.local/share}"/xorg/Xauthority
