# thinkmachine/dialog/.config/profile.d/dialog.sh
# XDG dir spec compliance
export DIALOGRC="${XDG_CONFIG_HOME:-$HOME/.config}"/dialog/dialogrc 
