# thinkmachine/nvidia/.config/profile.d/nvidia.sh
# XDG dir spec compliance
export __GL_SHADER_DISK_CACHE_PATH="${XDG_CACHE_HOME:-$HOME/.cache}"/.nv/GLCache
