# thinkmachine/node/.config/profile.d/node-gyp.sh
# XDG dir spec compliance
export npm_config_devdir="${XDG_CACHE_HOME:-$HOME/.cache}"/node-gyp
