# thinkmachine/node/.config/profile.d/node.sh
# XDG dir spec compliance
export NODE_REPL_HISTORY="${XDG_CACHE_HOME:-$HOME/.cache}"/node/repl_history
